# export_examples/navigator/APT39_G0087.json
{"description": "Enterprise techniques used by APT39, ATT&CK group G0087 v1.0", "name": "APT39 (G0087)", "domain": "mitre-enterprise", "version": "3.0", "techniques": [{"score": 1, "techniqueID": "T1546.010", "techniqueName": "AppInit DLLs", "comment": "[APT39](https://attack.mitre.org/groups/G0087) has used malware to set <code>LoadAppInit_DLLs</code> in the Registry key <code>SOFTWARE\\Microsoft\\Windows NT\\CurrentVersion\\Windows</code> in order to establish persistence.(Citation: FBI FLASH APT39 September 2020)"}, {"score": 1, "techniqueID": "T1140", "techniqueName": "Deobfuscate/Decode Files or Information", "comment": "[APT39](https://attack.mitre.org/groups/G0087) has used malware to decrypt encrypted CAB files.(Citation: FBI FLASH APT39 September 2020)"}, {"score": 1, "techniqueID": "T1553.006", "techniqueName": "Code Signing Policy Modification", "comment": "[APT39](https://attack.mitre.org/groups/G0087) has used malware to turn off the <code>RequireSigned</code> feature which ensures only signed DLLs can be run on Windows.(Citation: FBI FLASH APT39 September 2020)"}, {"score": 1, "techniqueID": "T1070.004", "techniqueName": "File Deletion", "comment": "[APT39](https://attack.mitre.org/groups/G0087) has used malware to delete files after they are deployed on a compromised host.(Citation: FBI FLASH APT39 September 2020)"}, {"score": 1, "techniqueID": "T1027", "techniqueName": "Obfuscated Files or Information", "comment": "[APT39](https://attack.mitre.org/groups/G0087) has used malware to drop encrypted CAB files.(Citation: FBI FLASH APT39 September 2020)"}, {"score": 1, "techniqueID": "T1041", "techniqueName": "Exfiltration Over C2 Channel", "comment": "[APT39](https://attack.mitre.org/groups/G0087) has exfiltrated stolen victim data through C2 communications.(Citation: FBI FLASH APT39 September 2020)"}, {"score": 1, "techniqueID": "T1074.001", "techniqueName": "Local Data Staging", "comment": "[APT39](https://attack.mitre.org/groups/G0087) has utilized tools to aggregate data prior to exfiltration.(Citation: FBI FLASH APT39 September 2020)"}, {"score": 1, "techniqueID": "T1197", "techniqueName": "BITS Jobs", "comment": "[APT39](https://attack.mitre.org/groups/G0087) has used the BITS protocol to exfiltrate stolen data from a compromised host.(Citation: FBI FLASH APT39 September 2020)"}, {"score": 1, "techniqueID": "T1083", "techniqueName": "File and Directory Discovery", "comment": "[APT39](https://attack.mitre.org/groups/G0087) has used tools with the ability to search for files on a compromised host.(Citation: FBI FLASH APT39 September 2020)"}, {"score": 1, "techniqueID": "T1012", "techniqueName": "Query Registry", "comment": "[APT39](https://attack.mitre.org/groups/G0087) has used various strains of malware to query the Registry.(Citation: FBI FLASH APT39 September 2020)"}, {"score": 1, "techniqueID": "T1056", "techniqueName": "Input Capture", "comment": "[APT39](https://attack.mitre.org/groups/G0087) has utilized tools to capture mouse movements.(Citation: FBI FLASH APT39 September 2020)"}, {"score": 1, "techniqueID": "T1059.005", "techniqueName": "Visual Basic", "comment": "[APT39](https://attack.mitre.org/groups/G0087) has utilized malicious VBS scripts in malware.(Citation: FBI FLASH APT39 September 2020)"}, {"score": 1, "techniqueID": "T1036.005", "techniqueName": "Match Legitimate Name or Location", "comment": "[APT39](https://attack.mitre.org/groups/G0087) has used malware disguised as Mozilla Firefox and a tool named mfevtpse.exe to proxy C2 communications, closely mimicking a legitimate McAfee file mfevtps.exe.(Citation: BitDefender Chafer May 2020)(Citation: FBI FLASH APT39 September 2020)"}, {"score": 1, "techniqueID": "T1190", "techniqueName": "Exploit Public-Facing Application", "comment": "[APT39](https://attack.mitre.org/groups/G0087) has used SQL injection for initial compromise.(Citation: Symantec Chafer February 2018)"}, {"score": 1, "techniqueID": "T1113", "techniqueName": "Screen Capture", "comment": "[APT39](https://attack.mitre.org/groups/G0087) has used a screen capture utility to take screenshots on a compromised host.(Citation: Symantec Chafer February 2018)(Citation: FBI FLASH APT39 September 2020)"}, {"score": 1, "techniqueID": "T1056.001", "techniqueName": "Keylogging", "comment": "[APT39](https://attack.mitre.org/groups/G0087) has used tools for capturing keystrokes.(Citation: Symantec Chafer February 2018)(Citation: FBI FLASH APT39 September 2020)"}, {"score": 1, "techniqueID": "T1005", "techniqueName": "Data from Local System", "comment": "[APT39](https://attack.mitre.org/groups/G0087) has used various tools to steal files from the compromised host.(Citation: Symantec Chafer February 2018)(Citation: FBI FLASH APT39 September 2020)"}, {"score": 1, "techniqueID": "T1115", "techniqueName": "Clipboard Data", "comment": "[APT39](https://attack.mitre.org/groups/G0087) has used tools capable of stealing contents of the clipboard.(Citation: Symantec Chafer February 2018)"}, {"score": 1, "techniqueID": "T1105", "techniqueName": "Ingress Tool Transfer", "comment": "[APT39](https://attack.mitre.org/groups/G0087) has downloaded tools to compromised hosts.(Citation: Symantec Chafer February 2018)(Citation: FBI FLASH APT39 September 2020)"}, {"score": 1, "techniqueID": "T1021.002", "techniqueName": "SMB/Windows Admin Shares", "comment": "[APT39](https://attack.mitre.org/groups/G0087) has used SMB for lateral movement.(Citation: Symantec Chafer February 2018)"}, {"score": 1, "techniqueID": "T1003", "techniqueName": "OS Credential Dumping", "comment": "[APT39](https://attack.mitre.org/groups/G0087) has used different versions of Mimikatz to obtain credentials.(Citation: BitDefender Chafer May 2020)"}, {"score": 1, "techniqueID": "T1071.001", "techniqueName": "Web Protocols", "comment": "[APT39](https://attack.mitre.org/groups/G0087) has used HTTP in communications with C2.(Citation: BitDefender Chafer May 2020)(Citation: FBI FLASH APT39 September 2020)"}, {"score": 1, "techniqueID": "T1059.006", "techniqueName": "Python", "comment": "[APT39](https://attack.mitre.org/groups/G0087) has used a command line utility and a network scanner written in python.(Citation: BitDefender Chafer May 2020)(Citation: FBI FLASH APT39 September 2020)"}, {"score": 1, "techniqueID": "T1071.004", "techniqueName": "DNS", "comment": "[APT39](https://attack.mitre.org/groups/G0087) has used remote access tools that leverage DNS in communications with C2.(Citation: BitDefender Chafer May 2020)"}, {"score": 1, "techniqueID": "T1102.002", "techniqueName": "Bidirectional Communication", "comment": "[APT39](https://attack.mitre.org/groups/G0087) has communicated with C2 through files uploaded to and downloaded from DropBox.(Citation: BitDefender Chafer May 2020)"}, {"score": 1, "techniqueID": "T1059.001", "techniqueName": "PowerShell", "comment": "[APT39](https://attack.mitre.org/groups/G0087) has used PowerShell to execute malicious code.(Citation: BitDefender Chafer May 2020)(Citation: Symantec Chafer February 2018)"}, {"score": 1, "techniqueID": "T1135", "techniqueName": "Network Share Discovery", "comment": "[APT39](https://attack.mitre.org/groups/G0087) has used the post exploitation tool [CrackMapExec](https://attack.mitre.org/software/S0488) to enumerate network shares.(Citation: BitDefender Chafer May 2020)"}, {"score": 1, "techniqueID": "T1569.002", "techniqueName": "Service Execution", "comment": "[APT39](https://attack.mitre.org/groups/G0087) has used post-exploitation tools including RemCom and the Non-sucking Service Manager (NSSM) to execute processes.(Citation: BitDefender Chafer May 2020)(Citation: Symantec Chafer February 2018)"}, {"score": 1, "techniqueID": "T1136.001", "techniqueName": "Local Account", "comment": "[APT39](https://attack.mitre.org/groups/G0087) has created accounts on multiple compromised hosts to perform actions within the network.(Citation: BitDefender Chafer May 2020)"}, {"score": 1, "techniqueID": "T1090.002", "techniqueName": "External Proxy", "comment": "[APT39](https://attack.mitre.org/groups/G0087) has used various tools to proxy C2 communications.(Citation: BitDefender Chafer May 2020)"}, {"score": 1, "techniqueID": "T1555", "techniqueName": "Credentials from Password Stores", "comment": "[APT39](https://attack.mitre.org/groups/G0087) has used the Smartftp Password Decryptor tool to decrypt FTP passwords.(Citation: BitDefender Chafer May 2020)"}, {"score": 1, "techniqueID": "T1110", "techniqueName": "Brute Force", "comment": "[APT39](https://attack.mitre.org/groups/G0087) has used Ncrack to reveal credentials.(Citation: FireEye APT39 Jan 2019)"}, {"score": 1, "techniqueID": "T1204.001", "techniqueName": "Malicious Link", "comment": "[APT39](https://attack.mitre.org/groups/G0087) has sent spearphishing emails in an attempt to lure users to click on a malicious link. (Citation: FireEye APT39 Jan 2019)(Citation: FBI FLASH APT39 September 2020)"}, {"score": 1, "techniqueID": "T1033", "techniqueName": "System Owner/User Discovery", "comment": "[APT39](https://attack.mitre.org/groups/G0087) used [Remexi](https://attack.mitre.org/software/S0375) to collect usernames from the system. (Citation: Symantec Chafer Dec 2015)"}, {"score": 1, "techniqueID": "T1027.002", "techniqueName": "Software Packing", "comment": "[APT39](https://attack.mitre.org/groups/G0087) has packed tools with UPX, and has repacked a modified version of [Mimikatz](https://attack.mitre.org/software/S0002) to thwart anti-virus detection.(Citation: FireEye APT39 Jan 2019)(Citation: BitDefender Chafer May 2020)"}, {"score": 1, "techniqueID": "T1204.002", "techniqueName": "Malicious File", "comment": "[APT39](https://attack.mitre.org/groups/G0087) has sent spearphishing emails in an attempt to lure users to click on a malicious attachment. (Citation: FireEye APT39 Jan 2019)(Citation: BitDefender Chafer May 2020)(Citation: Symantec Chafer February 2018)(Citation: FBI FLASH APT39 September 2020)"}, {"score": 1, "techniqueID": "T1505.003", "techniqueName": "Web Shell", "comment": "[APT39](https://attack.mitre.org/groups/G0087) has installed ANTAK and ASPXSPY web shells. (Citation: FireEye APT39 Jan 2019)"}, {"score": 1, "techniqueID": "T1021.004", "techniqueName": "SSH", "comment": "[APT39](https://attack.mitre.org/groups/G0087) used secure shell (SSH) to move laterally among their targets. (Citation: FireEye APT39 Jan 2019)"}, {"score": 1, "techniqueID": "T1078", "techniqueName": "Valid Accounts", "comment": "[APT39](https://attack.mitre.org/groups/G0087) has used stolen credentials to compromise Outlook Web Access (OWA). (Citation: FireEye APT39 Jan 2019)"}, {"score": 1, "techniqueID": "T1059", "techniqueName": "Command and Scripting Interpreter", "comment": "[APT39](https://attack.mitre.org/groups/G0087) has utilized AutoIt and custom scripts to perform internal reconnaissance. (Citation: FireEye APT39 Jan 2019)(Citation: FBI FLASH APT39 September 2020)"}, {"score": 1, "techniqueID": "T1046", "techniqueName": "Network Service Scanning", "comment": "[APT39](https://attack.mitre.org/groups/G0087) has used [CrackMapExec](https://attack.mitre.org/software/S0488) and a custom port scanner known as BLUETORCH for network scanning. (Citation: FireEye APT39 Jan 2019)(Citation: BitDefender Chafer May 2020)"}, {"score": 1, "techniqueID": "T1003.001", "techniqueName": "LSASS Memory", "comment": "[APT39](https://attack.mitre.org/groups/G0087) has used Mimikatz, Windows Credential Editor and ProcDump to dump credentials. (Citation: FireEye APT39 Jan 2019)"}, {"score": 1, "techniqueID": "T1090.001", "techniqueName": "Internal Proxy", "comment": "[APT39](https://attack.mitre.org/groups/G0087) used custom tools to create SOCK5 and custom protocol proxies between infected hosts. (Citation: FireEye APT39 Jan 2019)(Citation: BitDefender Chafer May 2020)"}, {"score": 1, "techniqueID": "T1560.001", "techniqueName": "Archive via Utility", "comment": "[APT39](https://attack.mitre.org/groups/G0087) has used WinRAR and 7-Zip to compress an archive stolen data. (Citation: FireEye APT39 Jan 2019)"}, {"score": 1, "techniqueID": "T1018", "techniqueName": "Remote System Discovery", "comment": "[APT39](https://attack.mitre.org/groups/G0087) has used [NBTscan](https://attack.mitre.org/software/S0590) and custom tools to discover remote systems. (Citation: FireEye APT39 Jan 2019)(Citation: BitDefender Chafer May 2020)(Citation: Symantec Chafer February 2018)"}, {"score": 1, "techniqueID": "T1547.009", "techniqueName": "Shortcut Modification", "comment": "[APT39](https://attack.mitre.org/groups/G0087) has modified LNK shortcuts. (Citation: FireEye APT39 Jan 2019)"}, {"score": 1, "techniqueID": "T1053.005", "techniqueName": "Scheduled Task", "comment": "[APT39](https://attack.mitre.org/groups/G0087) has created scheduled tasks for persistence. (Citation: FireEye APT39 Jan 2019)(Citation: BitDefender Chafer May 2020)(Citation: FBI FLASH APT39 September 2020)"}, {"score": 1, "techniqueID": "T1021.001", "techniqueName": "Remote Desktop Protocol", "comment": "[APT39](https://attack.mitre.org/groups/G0087) has been seen using RDP for lateral movement and persistence, in some cases employing the rdpwinst tool for mangement of multiple sessions. (Citation: FireEye APT39 Jan 2019)(Citation: BitDefender Chafer May 2020)"}, {"score": 1, "techniqueID": "T1547.001", "techniqueName": "Registry Run Keys / Startup Folder", "comment": "[APT39](https://attack.mitre.org/groups/G0087) has maintained persistence using the startup folder. (Citation: FireEye APT39 Jan 2019)"}, {"score": 1, "techniqueID": "T1566.002", "techniqueName": "Spearphishing Link", "comment": "[APT39](https://attack.mitre.org/groups/G0087) leveraged spearphishing emails with malicious links to initially compromise victims. (Citation: FireEye APT39 Jan 2019)(Citation: FBI FLASH APT39 September 2020)"}, {"score": 1, "techniqueID": "T1566.001", "techniqueName": "Spearphishing Attachment", "comment": "[APT39](https://attack.mitre.org/groups/G0087) leveraged spearphishing emails with malicious attachments to initially compromise victims. (Citation: FireEye APT39 Jan 2019)(Citation: Symantec Chafer February 2018)(Citation: FBI FLASH APT39 September 2020)"}], "gradient": {"colors": ["#ffffff", "#ff6666"], "minValue": 0, "maxValue": 1}, "legendItems": [{"label": "used by APT39", "color": "#ff6666"}]}