# export_examples/navigator/APT39_G0087.json
{"description": "Enterprise techniques used by APT39, ATT&CK group G0087 v1.0", "name": "APT39 (G0087)", "domain": "mitre-enterprise", "version": "2.2", "techniques": [{"score": 1, "techniqueID": "T1033", "techniqueName": "System Owner/User Discovery", "comment": "[APT39](https://attack.mitre.org/groups/G0087) used [Remexi](https://attack.mitre.org/software/S0375) to collect usernames from the system."}, {"score": 1, "techniqueID": "T1045", "techniqueName": "Software Packing", "comment": "[APT39](https://attack.mitre.org/groups/G0087) has repacked a modified version of [Mimikatz](https://attack.mitre.org/software/S0002) to thwart anti-virus detection."}, {"score": 1, "techniqueID": "T1204", "techniqueName": "User Execution", "comment": "[APT39](https://attack.mitre.org/groups/G0087) has sent spearphishing emails in an attempt to lure users to click on a malicious attachment or link."}, {"score": 1, "techniqueID": "T1100", "techniqueName": "Web Shell", "comment": "[APT39](https://attack.mitre.org/groups/G0087) has installed ANTAK and ASPXSPY web shells."}, {"score": 1, "techniqueID": "T1021", "techniqueName": "Remote Services", "comment": "[APT39](https://attack.mitre.org/groups/G0087) used secure shell (SSH) to move laterally among their targets."}, {"score": 1, "techniqueID": "T1078", "techniqueName": "Valid Accounts", "comment": "[APT39](https://attack.mitre.org/groups/G0087) has used stolen credentials to compromise Outlook Web Access (OWA)."}, {"score": 1, "techniqueID": "T1064", "techniqueName": "Scripting", "comment": "[APT39](https://attack.mitre.org/groups/G0087) utilized custom scripts to perform internal reconnaissance."}, {"score": 1, "techniqueID": "T1046", "techniqueName": "Network Service Scanning", "comment": "[APT39](https://attack.mitre.org/groups/G0087) used a custom port scanner known as BLUETORCH"}, {"score": 1, "techniqueID": "T1003", "techniqueName": "Credential Dumping", "comment": "[APT39](https://attack.mitre.org/groups/G0087) has used Mimikatz, Ncrack, Windows Credential Editor and ProcDump to dump credentials."}, {"score": 1, "techniqueID": "T1090", "techniqueName": "Connection Proxy", "comment": "[APT39](https://attack.mitre.org/groups/G0087) used custom tools to create SOCK5 proxies between infected hosts."}, {"score": 1, "techniqueID": "T1002", "techniqueName": "Data Compressed", "comment": "[APT39](https://attack.mitre.org/groups/G0087) has used WinRAR and 7-Zip to compress an archive stolen data."}, {"score": 1, "techniqueID": "T1016", "techniqueName": "System Network Configuration Discovery", "comment": "[APT39](https://attack.mitre.org/groups/G0087) has used nbtscan to discover vulnerable systems."}, {"score": 1, "techniqueID": "T1023", "techniqueName": "Shortcut Modification", "comment": "[APT39](https://attack.mitre.org/groups/G0087) has modified LNK shortcuts."}, {"score": 1, "techniqueID": "T1053", "techniqueName": "Scheduled Task", "comment": "[APT39](https://attack.mitre.org/groups/G0087) has created scheduled tasks."}, {"score": 1, "techniqueID": "T1076", "techniqueName": "Remote Desktop Protocol", "comment": "[APT39](https://attack.mitre.org/groups/G0087) has been seen using RDP for lateral movement and persistence."}, {"score": 1, "techniqueID": "T1060", "techniqueName": "Registry Run Keys / Startup Folder", "comment": "[APT39](https://attack.mitre.org/groups/G0087) has maintained persistence using the startup folder."}, {"score": 1, "techniqueID": "T1192", "techniqueName": "Spearphishing Link", "comment": "[APT39](https://attack.mitre.org/groups/G0087) leveraged spearphishing emails with malicious links to initially compromise victims."}, {"score": 1, "techniqueID": "T1193", "techniqueName": "Spearphishing Attachment", "comment": "[APT39](https://attack.mitre.org/groups/G0087) leveraged spearphishing emails with malicious attachments to initially compromise victims."}], "gradient": {"colors": ["#ffffff", "#ff6666"], "minValue": 0, "maxValue": 1}, "legendItems": [{"label": "used by APT39", "color": "#ff6666"}]}