# export_examples/navigator/APT32_G0050.json
{"description": "Enterprise techniques used by APT32, ATT&CK group G0050 v1.0", "name": "APT32 (G0050)", "domain": "mitre-enterprise", "version": "2.2", "techniques": [{"score": 1, "techniqueID": "T1143", "techniqueName": "Hidden Window", "comment": "[APT32](https://attack.mitre.org/groups/G0050) has used the WindowStyle parameter to conceal [PowerShell](https://attack.mitre.org/techniques/T1086) windows."}, {"score": 1, "techniqueID": "T1059", "techniqueName": "Command-Line Interface", "comment": "[APT32](https://attack.mitre.org/groups/G0050) has used cmd.exe for execution. "}, {"score": 1, "techniqueID": "T1045", "techniqueName": "Software Packing", "comment": "[APT32](https://attack.mitre.org/groups/G0050) uses UPX to pack their macOS backdoor."}, {"score": 1, "techniqueID": "T1222", "techniqueName": "File and Directory Permissions Modification", "comment": "[APT32](https://attack.mitre.org/groups/G0050)'s macOS backdoor changes the permission of the file it wants to execute to 755."}, {"score": 1, "techniqueID": "T1107", "techniqueName": "File Deletion", "comment": "[APT32](https://attack.mitre.org/groups/G0050)'s macOS backdoor can receive a \u201cdelete\u201d command."}, {"score": 1, "techniqueID": "T1158", "techniqueName": "Hidden Files and Directories", "comment": "[APT32](https://attack.mitre.org/groups/G0050)'s macOS backdoor hides the clientID file via a chflags function."}, {"score": 1, "techniqueID": "T1035", "techniqueName": "Service Execution", "comment": "[APT32](https://attack.mitre.org/groups/G0050)'s backdoor has used Windows services as a way to execute its malicious payload."}, {"score": 1, "techniqueID": "T1112", "techniqueName": "Modify Registry", "comment": "[APT32](https://attack.mitre.org/groups/G0050)'s backdoor has modified the Windows Registry to store the backdoor's configuration.\t\n"}, {"score": 1, "techniqueID": "T1002", "techniqueName": "Data Compressed", "comment": "[APT32](https://attack.mitre.org/groups/G0050)'s backdoor has used LZMA compression before exfiltration.\t"}, {"score": 1, "techniqueID": "T1012", "techniqueName": "Query Registry", "comment": "[APT32](https://attack.mitre.org/groups/G0050)'s backdoor can query the Windows Registry to gather system information."}, {"score": 1, "techniqueID": "T1083", "techniqueName": "File and Directory Discovery", "comment": "[APT32](https://attack.mitre.org/groups/G0050)'s backdoor possesses the capability to list files and directories on a machine.\t\n"}, {"score": 1, "techniqueID": "T1022", "techniqueName": "Data Encrypted", "comment": "[APT32](https://attack.mitre.org/groups/G0050) backdoors have encrypted data before exfiltration, including by using RC4 encryption."}, {"score": 1, "techniqueID": "T1041", "techniqueName": "Exfiltration Over Command and Control Channel", "comment": "[APT32](https://attack.mitre.org/groups/G0050)'s backdoor has exfiltrated data using the already opened channel with its C&C server."}, {"score": 1, "techniqueID": "T1203", "techniqueName": "Exploitation for Client Execution", "comment": "[APT32](https://attack.mitre.org/groups/G0050) has used RTF document that includes an exploit to execute malicious code. (CVE-2017-11882)"}, {"score": 1, "techniqueID": "T1065", "techniqueName": "Uncommonly Used Port", "comment": "[APT32](https://attack.mitre.org/groups/G0050) backdoor can use HTTP over an uncommon TCP port (e.g 14146) which is specified in the backdoor configuration."}, {"score": 1, "techniqueID": "T1043", "techniqueName": "Commonly Used Port", "comment": "[APT32](https://attack.mitre.org/groups/G0050) has used port 80 for C2 communications."}, {"score": 1, "techniqueID": "T1064", "techniqueName": "Scripting", "comment": "[APT32](https://attack.mitre.org/groups/G0050) has used macros, PowerShell scripts, COM scriptlets, and VBS scripts."}, {"score": 1, "techniqueID": "T1137", "techniqueName": "Office Application Startup", "comment": "[APT32](https://attack.mitre.org/groups/G0050) installed a backdoor macro in Microsoft Outlook for persistence."}, {"score": 1, "techniqueID": "T1003", "techniqueName": "Credential Dumping", "comment": "[APT32](https://attack.mitre.org/groups/G0050) used Mimikatz, GetPassword_x64, and  customized versions of Windows Credential Dumper, HookChangePassword, and Outlook Credential Dumper to harvest credentials."}, {"score": 1, "techniqueID": "T1060", "techniqueName": "Registry Run Keys / Startup Folder", "comment": "[APT32](https://attack.mitre.org/groups/G0050) established persistence using Registry Run keys, both to execute PowerShell and VBS scripts as well as to execute their backdoor directly."}, {"score": 1, "techniqueID": "T1192", "techniqueName": "Spearphishing Link", "comment": "[APT32](https://attack.mitre.org/groups/G0050) has sent spearphishing emails containing malicious links."}, {"score": 1, "techniqueID": "T1170", "techniqueName": "Mshta", "comment": "[APT32](https://attack.mitre.org/groups/G0050) has used mshta.exe for code execution."}, {"score": 1, "techniqueID": "T1087", "techniqueName": "Account Discovery", "comment": "[APT32](https://attack.mitre.org/groups/G0050) enumerated administrative users and DC servers using the commands <code>net localgroup administrators</code> and <code>net group \"Domain Controllers\" /domain</code>."}, {"score": 1, "techniqueID": "T1016", "techniqueName": "System Network Configuration Discovery", "comment": "[APT32](https://attack.mitre.org/groups/G0050) used the <code>ipconfig /all</code> command to gather the IP address from the system."}, {"score": 1, "techniqueID": "T1049", "techniqueName": "System Network Connections Discovery", "comment": "[APT32](https://attack.mitre.org/groups/G0050) used the <code>netstat -anpo tcp</code> command to display TCP connections on the victim's machine."}, {"score": 1, "techniqueID": "T1018", "techniqueName": "Remote System Discovery", "comment": "[APT32](https://attack.mitre.org/groups/G0050) used the <code>net view</code> command to show all shares available, including the administrative shares such as <code>C$</code> and <code>ADMIN$</code>. [APT32](https://attack.mitre.org/groups/G0050) also used the <code>ping</code> command."}, {"score": 1, "techniqueID": "T1096", "techniqueName": "NTFS File Attributes", "comment": "[APT32](https://attack.mitre.org/groups/G0050) used NTFS alternate data stream to hide their payloads."}, {"score": 1, "techniqueID": "T1047", "techniqueName": "Windows Management Instrumentation", "comment": "[APT32](https://attack.mitre.org/groups/G0050) used WMI to deploy their tools on remote machines and to gather information about the Outlook process."}, {"score": 1, "techniqueID": "T1031", "techniqueName": "Modify Existing Service", "comment": "[APT32](https://attack.mitre.org/groups/G0050) modified Windows Services to ensure PowerShell scripts were loaded on the system."}, {"score": 1, "techniqueID": "T1075", "techniqueName": "Pass the Hash", "comment": "[APT32](https://attack.mitre.org/groups/G0050) has used pass the hash for lateral movement."}, {"score": 1, "techniqueID": "T1046", "techniqueName": "Network Service Scanning", "comment": "[APT32](https://attack.mitre.org/groups/G0050) performed network scanning on the network to search for open ports, services, OS finger-printing, and other vulnerabilities."}, {"score": 1, "techniqueID": "T1097", "techniqueName": "Pass the Ticket", "comment": "[APT32](https://attack.mitre.org/groups/G0050) successfully gained remote access by using pass the ticket."}, {"score": 1, "techniqueID": "T1077", "techniqueName": "Windows Admin Shares", "comment": "[APT32](https://attack.mitre.org/groups/G0050) used [Net](https://attack.mitre.org/software/S0039) to use Windows' hidden network shares to copy their tools to remote machines for execution."}, {"score": 1, "techniqueID": "T1050", "techniqueName": "New Service", "comment": "[APT32](https://attack.mitre.org/groups/G0050) creates a Windows service to establish persistence."}, {"score": 1, "techniqueID": "T1189", "techniqueName": "Drive-by Compromise", "comment": "[APT32](https://attack.mitre.org/groups/G0050) has infected victims by tricking them into visiting compromised watering hole websites."}, {"score": 1, "techniqueID": "T1204", "techniqueName": "User Execution", "comment": "[APT32](https://attack.mitre.org/groups/G0050) has attempted to lure users to execute a malicious dropper delivered via a spearphishing attachment."}, {"score": 1, "techniqueID": "T1073", "techniqueName": "DLL Side-Loading", "comment": "[APT32](https://attack.mitre.org/groups/G0050) ran legitimately-signed executables from Symantec and McAfee which load a malicious DLL. The group also side-loads its backdoor by dropping a library and a legitimate, signed executable (AcroTranscoder)."}, {"score": 1, "techniqueID": "T1070", "techniqueName": "Indicator Removal on Host", "comment": "[APT32](https://attack.mitre.org/groups/G0050) has cleared select event log entries."}, {"score": 1, "techniqueID": "T1033", "techniqueName": "System Owner/User Discovery", "comment": "[APT32](https://attack.mitre.org/groups/G0050) collected the victim's username and executed the <code>whoami</code> command on the victim's machine."}, {"score": 1, "techniqueID": "T1193", "techniqueName": "Spearphishing Attachment", "comment": "[APT32](https://attack.mitre.org/groups/G0050) has sent spearphishing emails with a malicious executable disguised as a document or spreadsheet."}, {"score": 1, "techniqueID": "T1009", "techniqueName": "Binary Padding", "comment": "[APT32](https://attack.mitre.org/groups/G0050) includes garbage code to mislead anti-malware software and researchers."}, {"score": 1, "techniqueID": "T1082", "techniqueName": "System Information Discovery", "comment": "[APT32](https://attack.mitre.org/groups/G0050) has collected the OS version and computer name from victims. One of the group's backdoors can also query the Windows Registry to gather system information, and another macOS backdoor performs a fingerprint of the machine on its first connection to the C&C server."}, {"score": 1, "techniqueID": "T1216", "techniqueName": "Signed Script Proxy Execution", "comment": "[APT32](https://attack.mitre.org/groups/G0050) has used PubPrn.vbs within execution scripts to execute malware, possibly bypassing defenses."}, {"score": 1, "techniqueID": "T1100", "techniqueName": "Web Shell", "comment": "[APT32](https://attack.mitre.org/groups/G0050) has used Web shells to maintain access to victim websites."}, {"score": 1, "techniqueID": "T1105", "techniqueName": "Remote File Copy", "comment": "[APT32](https://attack.mitre.org/groups/G0050) has added JavaScript to victim websites to download additional frameworks that profile and compromise website visitors."}, {"score": 1, "techniqueID": "T1071", "techniqueName": "Standard Application Layer Protocol", "comment": "[APT32](https://attack.mitre.org/groups/G0050) has used JavaScript that communicates over HTTP or HTTPS to attacker controlled domains to download additional frameworks. The group has also used email for C2 via an Office macro. The group's backdoor can also exfiltrate data by encoding it in the subdomain field of DNS packets."}, {"score": 1, "techniqueID": "T1027", "techniqueName": "Obfuscated Files or Information", "comment": "[APT32](https://attack.mitre.org/groups/G0050) uses the Invoke-Obfuscation framework to obfuscate their PowerShell and also performs other code obfuscation. [APT32](https://attack.mitre.org/groups/G0050) has also encoded payloads using Base64 and a framework called \"Dont-Kill-My-Cat (DKMC). [APT32](https://attack.mitre.org/groups/G0050) also encrypts the library used for network exfiltration with AES-256 in CBC mode in their macOS backdoor."}, {"score": 1, "techniqueID": "T1053", "techniqueName": "Scheduled Task", "comment": "[APT32](https://attack.mitre.org/groups/G0050) has used scheduled tasks to persist on victim systems."}, {"score": 1, "techniqueID": "T1117", "techniqueName": "Regsvr32", "comment": "[APT32](https://attack.mitre.org/groups/G0050) created a [Scheduled Task](https://attack.mitre.org/techniques/T1053) that used regsvr32.exe to execute a COM scriptlet that dynamically downloaded a backdoor and injected it into memory. The group has also used regsvr32 to run their backdoor."}, {"score": 1, "techniqueID": "T1086", "techniqueName": "PowerShell", "comment": "[APT32](https://attack.mitre.org/groups/G0050) has used PowerShell-based tools, PowerShell one-liners, and shellcode loaders for execution."}, {"score": 1, "techniqueID": "T1094", "techniqueName": "Custom Command and Control Protocol", "comment": "[APT32](https://attack.mitre.org/groups/G0050) uses [Cobalt Strike](https://attack.mitre.org/software/S0154)'s malleable C2 functionality to blend in with network traffic. The group's backdoor can also exfiltrate data by encoding it in the subdomain field of DNS packets. Additionally, one of the group's macOS backdoors implements a specific format for the C2 packet involving random values."}, {"score": 1, "techniqueID": "T1068", "techniqueName": "Exploitation for Privilege Escalation", "comment": "[APT32](https://attack.mitre.org/groups/G0050) has used CVE-2016-7255 to escalate privileges."}, {"score": 1, "techniqueID": "T1017", "techniqueName": "Application Deployment Software", "comment": "[APT32](https://attack.mitre.org/groups/G0050) compromised McAfee ePO to move laterally by distributing malware as a software deployment task."}, {"score": 1, "techniqueID": "T1099", "techniqueName": "Timestomp", "comment": "[APT32](https://attack.mitre.org/groups/G0050) has used scheduled task raw XML with a backdated timestamp of June 2, 2016. The group has also set the creation time of the files dropped by the second stage of the exploit to match the creation time of kernel32.dll. Additionally, [APT32](https://attack.mitre.org/groups/G0050) has used a random value to modify the timestamp of the file storing the clientID."}, {"score": 1, "techniqueID": "T1078", "techniqueName": "Valid Accounts", "comment": "[APT32](https://attack.mitre.org/groups/G0050) has used legitimate local admin account credentials."}, {"score": 1, "techniqueID": "T1036", "techniqueName": "Masquerading", "comment": "[APT32](https://attack.mitre.org/groups/G0050) has used hidden or non-printing characters to help masquerade file names on a system, such as appending a Unicode no-break space character to a legitimate service name. They have also used by moving and renaming pubprn.vbs to a .txt file to avoid detection. Additionally, the group has renamed a NetCat binary to kb-10233.exe to masquerade as a Windows update."}], "gradient": {"colors": ["#ffffff", "#ff6666"], "minValue": 0, "maxValue": 1}, "legendItems": [{"label": "used by APT32", "color": "#ff6666"}]}