# export_examples/navigator/BRONZE BUTLER_G0060.json
{"description": "Enterprise techniques used by BRONZE BUTLER, ATT&CK group G0060 v1.0", "name": "BRONZE BUTLER (G0060)", "domain": "mitre-enterprise", "version": "2.2", "techniques": [{"score": 1, "techniqueID": "T1189", "techniqueName": "Drive-by Compromise", "comment": "[BRONZE BUTLER](https://attack.mitre.org/groups/G0060) compromised three Japanese websites using a Flash exploit to perform watering hole attacks."}, {"score": 1, "techniqueID": "T1203", "techniqueName": "Exploitation for Client Execution", "comment": "[BRONZE BUTLER](https://attack.mitre.org/groups/G0060) has exploited Microsoft Word vulnerability CVE-2014-4114 for execution."}, {"score": 1, "techniqueID": "T1204", "techniqueName": "User Execution", "comment": "[BRONZE BUTLER](https://attack.mitre.org/groups/G0060) has attempted to get users to launch malicious Microsoft Word attachments delivered via spearphishing emails."}, {"score": 1, "techniqueID": "T1193", "techniqueName": "Spearphishing Attachment", "comment": "[BRONZE BUTLER](https://attack.mitre.org/groups/G0060) used spearphishing emails with malicious Microsoft Word attachments to infect victims."}, {"score": 1, "techniqueID": "T1102", "techniqueName": "Web Service", "comment": "[BRONZE BUTLER](https://attack.mitre.org/groups/G0060)'s MSGET downloader uses a dead drop resolver to access malicious payloads."}, {"score": 1, "techniqueID": "T1002", "techniqueName": "Data Compressed", "comment": "[BRONZE BUTLER](https://attack.mitre.org/groups/G0060) has compressed data into password-protected RAR archives prior to exfiltration."}, {"score": 1, "techniqueID": "T1071", "techniqueName": "Standard Application Layer Protocol", "comment": "[BRONZE BUTLER](https://attack.mitre.org/groups/G0060) malware has used HTTP for C2."}, {"score": 1, "techniqueID": "T1105", "techniqueName": "Remote File Copy", "comment": "[BRONZE BUTLER](https://attack.mitre.org/groups/G0060) has used various tools to download files, including DGet (a similar tool to wget)."}, {"score": 1, "techniqueID": "T1036", "techniqueName": "Masquerading", "comment": "[BRONZE BUTLER](https://attack.mitre.org/groups/G0060) has given malware the same name as an existing file on the file share server to cause users to unwittingly launch and install the malware on additional systems."}, {"score": 1, "techniqueID": "T1088", "techniqueName": "Bypass User Account Control", "comment": "[BRONZE BUTLER](https://attack.mitre.org/groups/G0060) malware xxmm contains a UAC bypass tool for privilege escalation."}, {"score": 1, "techniqueID": "T1039", "techniqueName": "Data from Network Shared Drive", "comment": "[BRONZE BUTLER](https://attack.mitre.org/groups/G0060) has exfiltrated files stolen from file shares."}, {"score": 1, "techniqueID": "T1053", "techniqueName": "Scheduled Task", "comment": "[BRONZE BUTLER](https://attack.mitre.org/groups/G0060) has used [at](https://attack.mitre.org/software/S0110) and [schtasks](https://attack.mitre.org/software/S0111) to register a scheduled task to execute malware during lateral movement."}, {"score": 1, "techniqueID": "T1083", "techniqueName": "File and Directory Discovery", "comment": "[BRONZE BUTLER](https://attack.mitre.org/groups/G0060) has collected a list of files from the victim and uploaded it to its C2 server, and then created a new list of specific files to steal."}, {"score": 1, "techniqueID": "T1140", "techniqueName": "Deobfuscate/Decode Files or Information", "comment": "[BRONZE BUTLER](https://attack.mitre.org/groups/G0060) downloads encoded payloads and decodes them on the victim."}, {"score": 1, "techniqueID": "T1113", "techniqueName": "Screen Capture", "comment": "[BRONZE BUTLER](https://attack.mitre.org/groups/G0060) has used a tool to capture screenshots."}, {"score": 1, "techniqueID": "T1064", "techniqueName": "Scripting", "comment": "[BRONZE BUTLER](https://attack.mitre.org/groups/G0060) has used VBS, VBE, and batch scripts for execution."}, {"score": 1, "techniqueID": "T1022", "techniqueName": "Data Encrypted", "comment": "[BRONZE BUTLER](https://attack.mitre.org/groups/G0060) has compressed and encrypted data into password-protected RAR archives prior to exfiltration."}, {"score": 1, "techniqueID": "T1009", "techniqueName": "Binary Padding", "comment": "[BRONZE BUTLER](https://attack.mitre.org/groups/G0060) downloader code has included \"0\" characters at the end of the file to inflate the file size in a likely attempt to evade anti-virus detection."}, {"score": 1, "techniqueID": "T1032", "techniqueName": "Standard Cryptographic Protocol", "comment": "[BRONZE BUTLER](https://attack.mitre.org/groups/G0060) has used RC4 encryption (for Datper malware) and AES (for xxmm malware) to obfuscate HTTP traffic."}, {"score": 1, "techniqueID": "T1003", "techniqueName": "Credential Dumping", "comment": "[BRONZE BUTLER](https://attack.mitre.org/groups/G0060) has used various tools to perform credential dumping."}, {"score": 1, "techniqueID": "T1107", "techniqueName": "File Deletion", "comment": "The [BRONZE BUTLER](https://attack.mitre.org/groups/G0060) uploader or malware the uploader uses <code>command</code> to delete the RAR archives after they have been exfiltrated."}, {"score": 1, "techniqueID": "T1124", "techniqueName": "System Time Discovery", "comment": "[BRONZE BUTLER](https://attack.mitre.org/groups/G0060) has used <code>net time</code> to check the local time on a target system."}, {"score": 1, "techniqueID": "T1060", "techniqueName": "Registry Run Keys / Startup Folder", "comment": "[BRONZE BUTLER](https://attack.mitre.org/groups/G0060) has used a batch script that adds a Registry Run key to establish malware persistence."}, {"score": 1, "techniqueID": "T1097", "techniqueName": "Pass the Ticket", "comment": "[BRONZE BUTLER](https://attack.mitre.org/groups/G0060) has created forged Kerberos Ticket Granting Ticket (TGT) and Ticket Granting Service (TGS) tickets to maintain administrative access."}, {"score": 1, "techniqueID": "T1087", "techniqueName": "Account Discovery", "comment": "[BRONZE BUTLER](https://attack.mitre.org/groups/G0060) has used <code>net user /domain</code> to identify account information."}, {"score": 1, "techniqueID": "T1059", "techniqueName": "Command-Line Interface", "comment": "[BRONZE BUTLER](https://attack.mitre.org/groups/G0060) uses the command-line interface."}, {"score": 1, "techniqueID": "T1018", "techniqueName": "Remote System Discovery", "comment": "[BRONZE BUTLER](https://attack.mitre.org/groups/G0060) typically use <code>ping</code> and [Net](https://attack.mitre.org/software/S0039) to enumerate systems."}, {"score": 1, "techniqueID": "T1005", "techniqueName": "Data from Local System", "comment": "[BRONZE BUTLER](https://attack.mitre.org/groups/G0060) has exfiltrated files stolen from local systems."}, {"score": 1, "techniqueID": "T1024", "techniqueName": "Custom Cryptographic Protocol", "comment": "[BRONZE BUTLER](https://attack.mitre.org/groups/G0060) has used a tool called RarStar that encodes data with a custom XOR algorithm when posting it to a C2 server."}, {"score": 1, "techniqueID": "T1132", "techniqueName": "Data Encoding", "comment": "Several [BRONZE BUTLER](https://attack.mitre.org/groups/G0060) tools encode data with base64 when posting it to a C2 server."}, {"score": 1, "techniqueID": "T1086", "techniqueName": "PowerShell", "comment": "[BRONZE BUTLER](https://attack.mitre.org/groups/G0060) has used PowerShell for execution."}], "gradient": {"colors": ["#ffffff", "#ff6666"], "minValue": 0, "maxValue": 1}, "legendItems": [{"label": "used by BRONZE BUTLER", "color": "#ff6666"}]}