# export_examples/navigator/FIN8_G0061.json
{"description": "Enterprise techniques used by FIN8, ATT&CK group G0061 v1.0", "name": "FIN8 (G0061)", "domain": "mitre-enterprise", "version": "2.2", "techniques": [{"score": 1, "techniqueID": "T1027", "techniqueName": "Obfuscated Files or Information", "comment": "[FIN8](https://attack.mitre.org/groups/G0061) has used environment variables and standard input (stdin) to obfuscate command-line arguments. [FIN8](https://attack.mitre.org/groups/G0061) also obfuscates malicious macros delivered as payloads."}, {"score": 1, "techniqueID": "T1078", "techniqueName": "Valid Accounts", "comment": "[FIN8](https://attack.mitre.org/groups/G0061) has utilized [Valid Accounts](https://attack.mitre.org/techniques/T1078) during and [Persistence](https://attack.mitre.org/tactics/TA0003) and [Lateral Movement](https://attack.mitre.org/tactics/TA0008)."}, {"score": 1, "techniqueID": "T1105", "techniqueName": "Remote File Copy", "comment": "[FIN8](https://attack.mitre.org/groups/G0061) has used remote code execution to download subsequent payloads."}, {"score": 1, "techniqueID": "T1107", "techniqueName": "File Deletion", "comment": "[FIN8](https://attack.mitre.org/groups/G0061) has deleted tmp and prefetch files during post compromise cleanup activities."}, {"score": 1, "techniqueID": "T1048", "techniqueName": "Exfiltration Over Alternative Protocol", "comment": "[FIN8](https://attack.mitre.org/groups/G0061) has used FTP to exfiltrate collected data."}, {"score": 1, "techniqueID": "T1043", "techniqueName": "Commonly Used Port", "comment": "[FIN8](https://attack.mitre.org/groups/G0061) has tunneled RDP backdoors over port 443."}, {"score": 1, "techniqueID": "T1204", "techniqueName": "User Execution", "comment": "[FIN8](https://attack.mitre.org/groups/G0061) has leveraged both Spearphishing Link and Spearphishing Attachment attempting to gain User Execution."}, {"score": 1, "techniqueID": "T1002", "techniqueName": "Data Compressed", "comment": "[FIN8](https://attack.mitre.org/groups/G0061) has used RAR to compress collected data before [Exfiltration](https://attack.mitre.org/tactics/TA0010)."}, {"score": 1, "techniqueID": "T1003", "techniqueName": "Credential Dumping", "comment": "[FIN8](https://attack.mitre.org/groups/G0061) harvests credentials using Invoke-Mimikatz or Windows Credentials Editor (WCE)."}, {"score": 1, "techniqueID": "T1193", "techniqueName": "Spearphishing Attachment", "comment": "[FIN8](https://attack.mitre.org/groups/G0061) has distributed targeted emails containing Word documents with embedded malicious macros."}, {"score": 1, "techniqueID": "T1053", "techniqueName": "Scheduled Task", "comment": "[FIN8](https://attack.mitre.org/groups/G0061) has used scheduled tasks to maintain RDP backdoors."}, {"score": 1, "techniqueID": "T1018", "techniqueName": "Remote System Discovery", "comment": "[FIN8](https://attack.mitre.org/groups/G0061) uses [dsquery](https://attack.mitre.org/software/S0105) and other Active Directory utilities to enumerate hosts."}, {"score": 1, "techniqueID": "T1032", "techniqueName": "Standard Cryptographic Protocol", "comment": "[FIN8](https://attack.mitre.org/groups/G0061) has used the Plink utility to tunnel RDP back to C2 infrastructure."}, {"score": 1, "techniqueID": "T1192", "techniqueName": "Spearphishing Link", "comment": "[FIN8](https://attack.mitre.org/groups/G0061) has distributed targeted emails containing links to malicious documents with embedded macros."}, {"score": 1, "techniqueID": "T1086", "techniqueName": "PowerShell", "comment": "[FIN8](https://attack.mitre.org/groups/G0061)'s malicious spearphishing payloads are executed as [PowerShell](https://attack.mitre.org/techniques/T1086). [FIN8](https://attack.mitre.org/groups/G0061) has also used [PowerShell](https://attack.mitre.org/techniques/T1086) during [Lateral Movement](https://attack.mitre.org/tactics/TA0008) and [Credential Access](https://attack.mitre.org/tactics/TA0006)."}, {"score": 1, "techniqueID": "T1074", "techniqueName": "Data Staged", "comment": "[FIN8](https://attack.mitre.org/groups/G0061) aggregates staged data from a network into a single location."}, {"score": 1, "techniqueID": "T1063", "techniqueName": "Security Software Discovery", "comment": "[FIN8](https://attack.mitre.org/groups/G0061) has used Registry keys to detect and avoid executing in potential sandboxes."}, {"score": 1, "techniqueID": "T1064", "techniqueName": "Scripting", "comment": "[FIN8](https://attack.mitre.org/groups/G0061) has used a Batch file to automate frequently executed post compromise cleanup activities."}, {"score": 1, "techniqueID": "T1076", "techniqueName": "Remote Desktop Protocol", "comment": "[FIN8](https://attack.mitre.org/groups/G0061) has used RDP for [Lateral Movement](https://attack.mitre.org/tactics/TA0008)."}, {"score": 1, "techniqueID": "T1112", "techniqueName": "Modify Registry", "comment": "[FIN8](https://attack.mitre.org/groups/G0061) has deleted Registry keys during post compromise cleanup activities."}, {"score": 1, "techniqueID": "T1059", "techniqueName": "Command-Line Interface", "comment": "[FIN8](https://attack.mitre.org/groups/G0061) executes commands remotely via cmd.exe."}, {"score": 1, "techniqueID": "T1070", "techniqueName": "Indicator Removal on Host", "comment": "[FIN8](https://attack.mitre.org/groups/G0061) has cleared logs during post compromise cleanup activities."}, {"score": 1, "techniqueID": "T1047", "techniqueName": "Windows Management Instrumentation", "comment": "[FIN8](https://attack.mitre.org/groups/G0061)'s malicious spearphishing payloads use WMI to launch malware and spawn cmd.exe execution. [FIN8](https://attack.mitre.org/groups/G0061) has also used WMIC during and post compromise cleanup activities."}, {"score": 1, "techniqueID": "T1077", "techniqueName": "Windows Admin Shares", "comment": "[FIN8](https://attack.mitre.org/groups/G0061) has attempted to map to C$ on enumerated hosts to test the scope of their current credentials/context."}, {"score": 1, "techniqueID": "T1068", "techniqueName": "Exploitation for Privilege Escalation", "comment": "[FIN8](https://attack.mitre.org/groups/G0061) has exploited the CVE-2016-0167 local vulnerability."}], "gradient": {"colors": ["#ffffff", "#ff6666"], "minValue": 0, "maxValue": 1}, "legendItems": [{"label": "used by FIN8", "color": "#ff6666"}]}