# export_examples/navigator/Frankenstein_G0101.json
{"description": "Enterprise techniques used by Frankenstein, ATT&CK group G0101 v1.0", "name": "Frankenstein (G0101)", "domain": "mitre-enterprise", "version": "3.0", "techniques": [{"score": 1, "techniqueID": "T1140", "techniqueName": "Deobfuscate/Decode Files or Information", "comment": "[Frankenstein](https://attack.mitre.org/groups/G0101) has deobfuscated base64-encoded commands following the execution of a malicious script, which revealed a small script designed to obtain an additional payload.(Citation: Talos Frankenstein June 2019)  "}, {"score": 1, "techniqueID": "T1119", "techniqueName": "Automated Collection", "comment": "[Frankenstein](https://attack.mitre.org/groups/G0101) has enumerated hosts via [Empire](https://attack.mitre.org/software/S0363), gathering the username, domain name, machine name, and other system information.(Citation: Talos Frankenstein June 2019)"}, {"score": 1, "techniqueID": "T1005", "techniqueName": "Data from Local System", "comment": "[Frankenstein](https://attack.mitre.org/groups/G0101) has enumerated hosts via [Empire](https://attack.mitre.org/software/S0363), gathering various local system information.(Citation: Talos Frankenstein June 2019)"}, {"score": 1, "techniqueID": "T1033", "techniqueName": "System Owner/User Discovery", "comment": "[Frankenstein](https://attack.mitre.org/groups/G0101) has enumerated hosts, gathering username, machine name, and administrative permissions information.(Citation: Talos Frankenstein June 2019)"}, {"score": 1, "techniqueID": "T1573.001", "techniqueName": "Symmetric Cryptography", "comment": "[Frankenstein](https://attack.mitre.org/groups/G0101) has communicated with a C2 via an encrypted RC4 byte stream and AES-CBC.(Citation: Talos Frankenstein June 2019)"}, {"score": 1, "techniqueID": "T1082", "techniqueName": "System Information Discovery", "comment": "[Frankenstein](https://attack.mitre.org/groups/G0101) has enumerated hosts, looking for the system's machine name.(Citation: Talos Frankenstein June 2019)"}, {"score": 1, "techniqueID": "T1057", "techniqueName": "Process Discovery", "comment": "[Frankenstein](https://attack.mitre.org/groups/G0101) has enumerated hosts, looking to obtain a list of all currently running processes.(Citation: Talos Frankenstein June 2019)"}, {"score": 1, "techniqueID": "T1020", "techniqueName": "Automated Exfiltration", "comment": "[Frankenstein](https://attack.mitre.org/groups/G0101) has collected information via [Empire](https://attack.mitre.org/software/S0363), which is automatically sent the data back to the adversary's C2.(Citation: Talos Frankenstein June 2019)"}, {"score": 1, "techniqueID": "T1041", "techniqueName": "Exfiltration Over C2 Channel", "comment": "[Frankenstein](https://attack.mitre.org/groups/G0101) has collected information via [Empire](https://attack.mitre.org/software/S0363), which is automatically sent the data back to the adversary's C2.(Citation: Talos Frankenstein June 2019)"}, {"score": 1, "techniqueID": "T1105", "techniqueName": "Ingress Tool Transfer", "comment": "[Frankenstein](https://attack.mitre.org/groups/G0101) has uploaded and downloaded files to utilize additional plugins.(Citation: Talos Frankenstein June 2019)"}, {"score": 1, "techniqueID": "T1016", "techniqueName": "System Network Configuration Discovery", "comment": "[Frankenstein](https://attack.mitre.org/groups/G0101) has enumerated hosts, looking for the public IP address of the system.(Citation: Talos Frankenstein June 2019)"}, {"score": 1, "techniqueID": "T1003", "techniqueName": "OS Credential Dumping", "comment": "[Frankenstein](https://attack.mitre.org/groups/G0101) has harvested credentials from the victim's machine using [Empire](https://attack.mitre.org/software/S0363).(Citation: Talos Frankenstein June 2019)"}, {"score": 1, "techniqueID": "T1518.001", "techniqueName": "Security Software Discovery", "comment": "[Frankenstein](https://attack.mitre.org/groups/G0101) has used WMI queries to detect if virtualization environments or analysis tools were running on the system.(Citation: Talos Frankenstein June 2019)"}, {"score": 1, "techniqueID": "T1497.001", "techniqueName": "System Checks", "comment": "[Frankenstein](https://attack.mitre.org/groups/G0101) has used WMI queries to check if various security applications were running, including VMWare and Virtualbox.(Citation: Talos Frankenstein June 2019)"}, {"score": 1, "techniqueID": "T1047", "techniqueName": "Windows Management Instrumentation", "comment": "[Frankenstein](https://attack.mitre.org/groups/G0101) has used WMI queries to check if various security applications were running, as well as the operating system version.(Citation: Talos Frankenstein June 2019)"}, {"score": 1, "techniqueID": "T1127.001", "techniqueName": "MSBuild", "comment": "[Frankenstein](https://attack.mitre.org/groups/G0101) has used MSbuild to execute an actor-created file.(Citation: Talos Frankenstein June 2019)"}, {"score": 1, "techniqueID": "T1027", "techniqueName": "Obfuscated Files or Information", "comment": "[Frankenstein](https://attack.mitre.org/groups/G0101) has run encoded commands from the command line.(Citation: Talos Frankenstein June 2019)"}, {"score": 1, "techniqueID": "T1059.005", "techniqueName": "Visual Basic", "comment": "[Frankenstein](https://attack.mitre.org/groups/G0101) has used Word documents that prompts the victim to enable macros and run a Visual Basic script.(Citation: Talos Frankenstein June 2019)"}, {"score": 1, "techniqueID": "T1203", "techniqueName": "Exploitation for Client Execution", "comment": "[Frankenstein](https://attack.mitre.org/groups/G0101) has used CVE-2017-11882 to execute code on the victim's machine.(Citation: Talos Frankenstein June 2019)"}, {"score": 1, "techniqueID": "T1204.002", "techniqueName": "Malicious File", "comment": "[Frankenstein](https://attack.mitre.org/groups/G0101) has used trojanized Microsoft Word documents sent via email, which prompted the victim to enable macros.(Citation: Talos Frankenstein June 2019)"}, {"score": 1, "techniqueID": "T1059.003", "techniqueName": "Windows Command Shell", "comment": "[Frankenstein](https://attack.mitre.org/groups/G0101) has run a command script to set up persistence as a scheduled task named \"WinUpdate\", as well as other encoded commands from the command-line.(Citation: Talos Frankenstein June 2019)"}, {"score": 1, "techniqueID": "T1221", "techniqueName": "Template Injection", "comment": "[Frankenstein](https://attack.mitre.org/groups/G0101) has used trojanized documents that retrieve remote templates from an adversary-controlled website.(Citation: Talos Frankenstein June 2019)"}, {"score": 1, "techniqueID": "T1059.001", "techniqueName": "PowerShell", "comment": "[Frankenstein](https://attack.mitre.org/groups/G0101) has used PowerShell to run a series of base64-encoded commands, that acted as a stager and enumerated hosts.(Citation: Talos Frankenstein June 2019)\t"}, {"score": 1, "techniqueID": "T1053.005", "techniqueName": "Scheduled Task", "comment": "[Frankenstein](https://attack.mitre.org/groups/G0101) has established persistence through a scheduled task using the command: <code> /Create /F /SC DAILY /ST 09:00 /TN WinUpdate /TR </code>, named \"WinUpdate\".(Citation: Talos Frankenstein June 2019)\t"}, {"score": 1, "techniqueID": "T1566.001", "techniqueName": "Spearphishing Attachment", "comment": "[Frankenstein](https://attack.mitre.org/groups/G0101) has used spearphishing emails to send trojanized Microsoft Word documents.(Citation: Talos Frankenstein June 2019)  "}], "gradient": {"colors": ["#ffffff", "#ff6666"], "minValue": 0, "maxValue": 1}, "legendItems": [{"label": "used by Frankenstein", "color": "#ff6666"}]}