# export_examples/navigator/menuPass_G0045.json
{"description": "Enterprise techniques used by menuPass, ATT&CK group G0045 v1.0", "name": "menuPass (G0045)", "domain": "mitre-enterprise", "version": "2.2", "techniques": [{"score": 1, "techniqueID": "T1056", "techniqueName": "Input Capture", "comment": "[menuPass](https://attack.mitre.org/groups/G0045) has used key loggers to steal usernames and passwords."}, {"score": 1, "techniqueID": "T1005", "techniqueName": "Data from Local System", "comment": "[menuPass](https://attack.mitre.org/groups/G0045) has collected various files from the compromised computers.\n"}, {"score": 1, "techniqueID": "T1022", "techniqueName": "Data Encrypted", "comment": "[menuPass](https://attack.mitre.org/groups/G0045) has encrypted files and information before exfiltration."}, {"score": 1, "techniqueID": "T1036", "techniqueName": "Masquerading", "comment": "[menuPass](https://attack.mitre.org/groups/G0045) has been seen changing malicious files to appear legitimate. They have also renamed [certutil](https://attack.mitre.org/software/S0160) and move it to a different location on system to avoid detection based on use of the tool. The group has also used [esentutl](https://attack.mitre.org/software/S0404) to change file extensions to avoid detection."}, {"score": 1, "techniqueID": "T1027", "techniqueName": "Obfuscated Files or Information", "comment": "[menuPass](https://attack.mitre.org/groups/G0045) has encoded strings in its malware with base64 as well as with a simple, single-byte XOR obfuscation using key 0x40."}, {"score": 1, "techniqueID": "T1064", "techniqueName": "Scripting", "comment": "[menuPass](https://attack.mitre.org/groups/G0045) has used malicious macros embedded inside Office documents to execute files."}, {"score": 1, "techniqueID": "T1093", "techniqueName": "Process Hollowing", "comment": "[menuPass](https://attack.mitre.org/groups/G0045) has used process hollowing in iexplore.exe to load the [RedLeaves](https://attack.mitre.org/software/S0153) implant."}, {"score": 1, "techniqueID": "T1140", "techniqueName": "Deobfuscate/Decode Files or Information", "comment": "[menuPass](https://attack.mitre.org/groups/G0045) has used [certutil](https://attack.mitre.org/software/S0160) in a macro to decode base64-encoded content contained in a dropper document attached to an email. The group has also used <code>certutil -decode</code> to decode files on the victim\u2019s machine when dropping [UPPERCUT](https://attack.mitre.org/software/S0275)."}, {"score": 1, "techniqueID": "T1107", "techniqueName": "File Deletion", "comment": "A [menuPass](https://attack.mitre.org/groups/G0045) macro deletes files after it has decoded and decompressed them."}, {"score": 1, "techniqueID": "T1199", "techniqueName": "Trusted Relationship", "comment": "[menuPass](https://attack.mitre.org/groups/G0045) has used legitimate access granted to Managed Service Providers in order to access victims of interest."}, {"score": 1, "techniqueID": "T1204", "techniqueName": "User Execution", "comment": "[menuPass](https://attack.mitre.org/groups/G0045) has attempted to get victims to open malicious files such as Windows Shortcuts (.lnk) and/or Microsoft Office documents, sent via email as part of spearphishing campaigns."}, {"score": 1, "techniqueID": "T1193", "techniqueName": "Spearphishing Attachment", "comment": "[menuPass](https://attack.mitre.org/groups/G0045) has sent malicious Office documents via email as part of spearphishing campaigns as well as executables disguised as documents."}, {"score": 1, "techniqueID": "T1078", "techniqueName": "Valid Accounts", "comment": "[menuPass](https://attack.mitre.org/groups/G0045) has used valid accounts shared between Managed Service Providers and clients to move between the two environments."}, {"score": 1, "techniqueID": "T1016", "techniqueName": "System Network Configuration Discovery", "comment": "[menuPass](https://attack.mitre.org/groups/G0045) has used several tools to scan for open NetBIOS nameservers and enumerate NetBIOS sessions."}, {"score": 1, "techniqueID": "T1049", "techniqueName": "System Network Connections Discovery", "comment": "[menuPass](https://attack.mitre.org/groups/G0045) has used <code>net use</code> to conduct connectivity checks to machines."}, {"score": 1, "techniqueID": "T1047", "techniqueName": "Windows Management Instrumentation", "comment": "[menuPass](https://attack.mitre.org/groups/G0045) uses a modified version of pentesting script wmiexec.vbs, which logs into a remote machine using WMI."}, {"score": 1, "techniqueID": "T1076", "techniqueName": "Remote Desktop Protocol", "comment": "[menuPass](https://attack.mitre.org/groups/G0045) has used RDP connections to move across the victim network."}, {"score": 1, "techniqueID": "T1039", "techniqueName": "Data from Network Shared Drive", "comment": "[menuPass](https://attack.mitre.org/groups/G0045) has collected data from remote systems by mounting network shares with <code>net use</code> and using Robocopy to transfer data."}, {"score": 1, "techniqueID": "T1086", "techniqueName": "PowerShell", "comment": "[menuPass](https://attack.mitre.org/groups/G0045) uses [PowerSploit](https://attack.mitre.org/software/S0194) to inject shellcode into PowerShell."}, {"score": 1, "techniqueID": "T1002", "techniqueName": "Data Compressed", "comment": "[menuPass](https://attack.mitre.org/groups/G0045) has compressed files before exfiltration using TAR and RAR."}, {"score": 1, "techniqueID": "T1053", "techniqueName": "Scheduled Task", "comment": "[menuPass](https://attack.mitre.org/groups/G0045) has used a script (atexec.py) to execute a command on a target machine via Task Scheduler."}, {"score": 1, "techniqueID": "T1087", "techniqueName": "Account Discovery", "comment": "[menuPass](https://attack.mitre.org/groups/G0045) has used the Microsoft administration tool csvde.exe to export Active Directory data."}, {"score": 1, "techniqueID": "T1090", "techniqueName": "Connection Proxy", "comment": "[menuPass](https://attack.mitre.org/groups/G0045) has used a global service provider's IP as a proxy for C2 traffic from a victim."}, {"score": 1, "techniqueID": "T1059", "techniqueName": "Command-Line Interface", "comment": "[menuPass](https://attack.mitre.org/groups/G0045) executes commands using a command-line interface and reverse shell. The group has used a modified version of pentesting script wmiexec.vbs to execute commands."}, {"score": 1, "techniqueID": "T1003", "techniqueName": "Credential Dumping", "comment": "[menuPass](https://attack.mitre.org/groups/G0045) has used a modified version of pentesting tools wmiexec.vbs and secretsdump.py to dump credentials."}, {"score": 1, "techniqueID": "T1074", "techniqueName": "Data Staged", "comment": "[menuPass](https://attack.mitre.org/groups/G0045) stages data prior to exfiltration in multi-part archives, often saved in the Recycle Bin."}, {"score": 1, "techniqueID": "T1021", "techniqueName": "Remote Services", "comment": "[menuPass](https://attack.mitre.org/groups/G0045) has used Putty Secure Copy Client (PSCP) to transfer data."}, {"score": 1, "techniqueID": "T1073", "techniqueName": "DLL Side-Loading", "comment": "[menuPass](https://attack.mitre.org/groups/G0045) has used DLL side-loading to launch versions of Mimikatz and PwDump6 as well as [UPPERCUT](https://attack.mitre.org/software/S0275)."}, {"score": 1, "techniqueID": "T1038", "techniqueName": "DLL Search Order Hijacking", "comment": "[menuPass](https://attack.mitre.org/groups/G0045) has used DLL search order hijacking."}, {"score": 1, "techniqueID": "T1046", "techniqueName": "Network Service Scanning", "comment": "[menuPass](https://attack.mitre.org/groups/G0045) has used tcping.exe, similar to [Ping](https://attack.mitre.org/software/S0097), to probe port status on systems of interest."}, {"score": 1, "techniqueID": "T1018", "techniqueName": "Remote System Discovery", "comment": "[menuPass](https://attack.mitre.org/groups/G0045) uses scripts to enumerate IP ranges on the victim network. [menuPass](https://attack.mitre.org/groups/G0045) has also issued the command <code>net view /domain</code> to a [PlugX](https://attack.mitre.org/software/S0013) implant to gather information about remote systems on the network."}, {"score": 1, "techniqueID": "T1105", "techniqueName": "Remote File Copy", "comment": "[menuPass](https://attack.mitre.org/groups/G0045) has installed updates and new malware on victims."}], "gradient": {"colors": ["#ffffff", "#ff6666"], "minValue": 0, "maxValue": 1}, "legendItems": [{"label": "used by menuPass", "color": "#ff6666"}]}