# export_examples/navigator/TEMP.Veles_G0088.json
{"description": "Enterprise techniques used by TEMP.Veles, ATT&CK group G0088 v1.0", "name": "TEMP.Veles (G0088)", "domain": "mitre-enterprise", "version": "3.0", "techniques": [{"score": 1, "techniqueID": "T1583.003", "techniqueName": "Virtual Private Server", "comment": "[TEMP.Veles](https://attack.mitre.org/groups/G0088) has used Virtual Private Server (VPS) infrastructure.(Citation: FireEye TRITON 2019)"}, {"score": 1, "techniqueID": "T1571", "techniqueName": "Non-Standard Port", "comment": "[TEMP.Veles](https://attack.mitre.org/groups/G0088) has used port-protocol mismatches on ports such as 443, 4444, 8531, and 50501 during C2.(Citation: FireEye TRITON 2019)"}, {"score": 1, "techniqueID": "T1027.005", "techniqueName": "Indicator Removal from Tools", "comment": "[TEMP.Veles](https://attack.mitre.org/groups/G0088) has modified files based on the open-source project cryptcat in an apparent attempt to decrease AV detection rates.(Citation: FireEye TEMP.Veles 2018)"}, {"score": 1, "techniqueID": "T1059.001", "techniqueName": "PowerShell", "comment": "[TEMP.Veles](https://attack.mitre.org/groups/G0088) has used a publicly-available PowerShell-based tool, WMImplant.(Citation: FireEye TEMP.Veles 2018) The group has also used PowerShell to perform [Timestomp](https://attack.mitre.org/techniques/T1070/006)ing.(Citation: FireEye TRITON 2019)"}, {"score": 1, "techniqueID": "T1043", "techniqueName": "Commonly Used Port", "comment": "[TEMP.Veles](https://attack.mitre.org/groups/G0088) has used port 443 for C2.(Citation: FireEye TRITON 2019)"}, {"score": 1, "techniqueID": "T1065", "techniqueName": "Uncommonly Used Port", "comment": "[TEMP.Veles](https://attack.mitre.org/groups/G0088) has used ports 4444, 8531, and 50501 for C2.(Citation: FireEye TRITON 2019)"}, {"score": 1, "techniqueID": "T1505.003", "techniqueName": "Web Shell", "comment": "[TEMP.Veles](https://attack.mitre.org/groups/G0088) has planted Web shells on Outlook Exchange servers.(Citation: FireEye TRITON 2019)"}, {"score": 1, "techniqueID": "T1133", "techniqueName": "External Remote Services", "comment": "[TEMP.Veles](https://attack.mitre.org/groups/G0088) has used a VPN to persist in the victim environment.(Citation: FireEye TRITON 2019)"}, {"score": 1, "techniqueID": "T1021.001", "techniqueName": "Remote Desktop Protocol", "comment": "[TEMP.Veles](https://attack.mitre.org/groups/G0088) utilized RDP throughout an operation.(Citation: FireEye TRITON 2019)\t"}, {"score": 1, "techniqueID": "T1078", "techniqueName": "Valid Accounts", "comment": "[TEMP.Veles](https://attack.mitre.org/groups/G0088) has used compromised VPN accounts.(Citation: FireEye TRITON 2019)\t"}, {"score": 1, "techniqueID": "T1021.004", "techniqueName": "SSH", "comment": "[TEMP.Veles](https://attack.mitre.org/groups/G0088) has relied on encrypted SSH-based tunnels to transfer tools and for remote command/program execution.(Citation: FireEye TRITON 2019)"}, {"score": 1, "techniqueID": "T1074.001", "techniqueName": "Local Data Staging", "comment": "[TEMP.Veles](https://attack.mitre.org/groups/G0088) has created staging folders in directories that were infrequently used by legitimate users or processes.(Citation: FireEye TRITON 2019)"}, {"score": 1, "techniqueID": "T1070.004", "techniqueName": "File Deletion", "comment": "[TEMP.Veles](https://attack.mitre.org/groups/G0088) routinely deleted tools, logs, and other files after they were finished with them.(Citation: FireEye TRITON 2019)"}, {"score": 1, "techniqueID": "T1070.006", "techniqueName": "Timestomp", "comment": "[TEMP.Veles](https://attack.mitre.org/groups/G0088) used timestomping to modify the $STANDARD_INFORMATION attribute on tools.(Citation: FireEye TRITON 2019)"}, {"score": 1, "techniqueID": "T1053.005", "techniqueName": "Scheduled Task", "comment": "[TEMP.Veles](https://attack.mitre.org/groups/G0088) has used scheduled task XML triggers.(Citation: FireEye TRITON 2019)"}, {"score": 1, "techniqueID": "T1546.012", "techniqueName": "Image File Execution Options Injection", "comment": "[TEMP.Veles](https://attack.mitre.org/groups/G0088) has modified and added entries within <code>HKEY_LOCAL_MACHINE\\SOFTWARE\\Microsoft\\Windows NT\\CurrentVersion\\Image File Execution Options</code> to maintain persistence.(Citation: FireEye TRITON 2019)\t\n"}, {"score": 1, "techniqueID": "T1003.001", "techniqueName": "LSASS Memory", "comment": "[TEMP.Veles](https://attack.mitre.org/groups/G0088) has used [Mimikatz](https://attack.mitre.org/software/S0002) and a custom tool, SecHack, to harvest credentials. (Citation: FireEye TRITON 2019)"}, {"score": 1, "techniqueID": "T1036.005", "techniqueName": "Match Legitimate Name or Location", "comment": "[TEMP.Veles](https://attack.mitre.org/groups/G0088) has renamed files to look like legitimate files, such as Windows update files or Schneider Electric application files.(Citation: FireEye TRITON 2019)\n"}, {"score": 1, "techniqueID": "T1311", "techniqueName": "Dynamic DNS", "comment": "[TEMP.Veles](https://attack.mitre.org/groups/G0088) has used dynamic DNS.(Citation: FireEye TRITON 2019)"}, {"score": 1, "techniqueID": "T1329", "techniqueName": "Acquire and/or use 3rd party infrastructure services", "comment": "[TEMP.Veles](https://attack.mitre.org/groups/G0088) has used Virtual Private Server (VPS) infrastructure.(Citation: FireEye TRITON 2019)"}, {"score": 1, "techniqueID": "T0886", "techniqueName": "Remote Services", "comment": "[XENOTIME](https://collaborate.mitre.org/attackics/index.php/Group/G0001) utilized remote desktop protocol (RDP) jump boxes to move into the ICS environment.(Citation: Dragos TRISIS December 2017"}, {"score": 1, "techniqueID": "T0822", "techniqueName": "External Remote Services", "comment": "[XENOTIME](https://collaborate.mitre.org/attackics/index.php/Group/G0001) utilized remote desktop protocol (RDP) jump boxes to move into the ICS environment.(Citation: Dragos TRISIS Dec 2017)"}, {"score": 1, "techniqueID": "T0859", "techniqueName": "Valid Accounts", "comment": "[XENOTIME](https://collaborate.mitre.org/attackics/index.php/Group/G0001) used valid credentials when laterally moving through RDP jump boxes into the ICS environment.(Citation: Dragos TRISIS Dec 2017)"}, {"score": 1, "techniqueID": "T0817", "techniqueName": "Drive-by Compromise", "comment": "[XENOTIME](https://collaborate.mitre.org/attackics/index.php/Group/G0001) utilizes watering hole websites to target industrial employees.(Citation: CyberScoop Xenotime May 2018)"}, {"score": 1, "techniqueID": "T0862", "techniqueName": "Supply Chain Compromise", "comment": "[XENOTIME](https://collaborate.mitre.org/attackics/index.php/Group/G0001) targeted several ICS vendors and manufacturers.(Citation: Dragos Aug 2019)"}], "gradient": {"colors": ["#ffffff", "#ff6666"], "minValue": 0, "maxValue": 1}, "legendItems": [{"label": "used by TEMP.Veles", "color": "#ff6666"}]}