# export_examples/navigator/TEMP.Veles_G0088.json
{"description": "Enterprise techniques used by TEMP.Veles, ATT&CK group G0088 v1.0", "name": "TEMP.Veles (G0088)", "domain": "mitre-enterprise", "version": "2.2", "techniques": [{"score": 1, "techniqueID": "T1066", "techniqueName": "Indicator Removal from Tools", "comment": "[TEMP.Veles](https://attack.mitre.org/groups/G0088) has modified files based on the open-source project cryptcat in an apparent attempt to decrease AV detection rates."}, {"score": 1, "techniqueID": "T1086", "techniqueName": "PowerShell", "comment": "[TEMP.Veles](https://attack.mitre.org/groups/G0088) has used a publicly-available PowerShell-based tool, WMImplant. The group has also used PowerShell to perform [Timestomp](https://attack.mitre.org/techniques/T1099)ing."}, {"score": 1, "techniqueID": "T1043", "techniqueName": "Commonly Used Port", "comment": "[TEMP.Veles](https://attack.mitre.org/groups/G0088) has used port 443 for C2."}, {"score": 1, "techniqueID": "T1065", "techniqueName": "Uncommonly Used Port", "comment": "[TEMP.Veles](https://attack.mitre.org/groups/G0088) has used ports 4444, 8531, and 50501 for C2."}, {"score": 1, "techniqueID": "T1100", "techniqueName": "Web Shell", "comment": "[TEMP.Veles](https://attack.mitre.org/groups/G0088) has planted webshells on Outlook Exchange servers."}, {"score": 1, "techniqueID": "T1133", "techniqueName": "External Remote Services", "comment": "[TEMP.Veles](https://attack.mitre.org/groups/G0088) has used a VPN to persist in the victim environment."}, {"score": 1, "techniqueID": "T1076", "techniqueName": "Remote Desktop Protocol", "comment": "[TEMP.Veles](https://attack.mitre.org/groups/G0088) utilized RDP throughout an operation.\t"}, {"score": 1, "techniqueID": "T1078", "techniqueName": "Valid Accounts", "comment": "[TEMP.Veles](https://attack.mitre.org/groups/G0088) has used compromised VPN accounts.\t"}, {"score": 1, "techniqueID": "T1021", "techniqueName": "Remote Services", "comment": "[TEMP.Veles](https://attack.mitre.org/groups/G0088) has relied on encrypted SSH-based tunnels to transfer tools and for remote command/program execution."}, {"score": 1, "techniqueID": "T1074", "techniqueName": "Data Staged", "comment": "[TEMP.Veles](https://attack.mitre.org/groups/G0088) has created staging folders in directories that were infrequently used by legitimate users or processes."}, {"score": 1, "techniqueID": "T1107", "techniqueName": "File Deletion", "comment": "[TEMP.Veles](https://attack.mitre.org/groups/G0088) routinely deleted tools, logs, and other files after they were finished with them."}, {"score": 1, "techniqueID": "T1099", "techniqueName": "Timestomp", "comment": "[TEMP.Veles](https://attack.mitre.org/groups/G0088) used timestomping to modify the $STANDARD_INFORMATION attribute on tools."}, {"score": 1, "techniqueID": "T1053", "techniqueName": "Scheduled Task", "comment": "[TEMP.Veles](https://attack.mitre.org/groups/G0088) has used scheduled task XML triggers."}, {"score": 1, "techniqueID": "T1183", "techniqueName": "Image File Execution Options Injection", "comment": "[TEMP.Veles](https://attack.mitre.org/groups/G0088) has modified and added entries within <code>HKEY_LOCAL_MACHINE\\SOFTWARE\\Microsoft\\Windows NT\\CurrentVersion\\Image File Execution Options</code> to maintain persistence.\t\n"}, {"score": 1, "techniqueID": "T1003", "techniqueName": "Credential Dumping", "comment": "[TEMP.Veles](https://attack.mitre.org/groups/G0088) has used [Mimikatz](https://attack.mitre.org/software/S0002) and a custom tool, SecHack, to harvest credentials."}, {"score": 1, "techniqueID": "T1036", "techniqueName": "Masquerading", "comment": "[TEMP.Veles](https://attack.mitre.org/groups/G0088) has performed a variety of methods to look like valid users, including renaming files and mimicking legitimate administrator activities.\n"}, {"score": 1, "techniqueID": "T1311", "techniqueName": "Dynamic DNS", "comment": "[TEMP.Veles](https://attack.mitre.org/groups/G0088) has used dynamic DNS."}, {"score": 1, "techniqueID": "T1329", "techniqueName": "Acquire and/or use 3rd party infrastructure services", "comment": "[TEMP.Veles](https://attack.mitre.org/groups/G0088) has used Virtual Private Server (VPS) infrastructure."}], "gradient": {"colors": ["#ffffff", "#ff6666"], "minValue": 0, "maxValue": 1}, "legendItems": [{"label": "used by TEMP.Veles", "color": "#ff6666"}]}