# export_examples/navigator/Threat Group-3390_G0027.json
{"description": "Enterprise techniques used by Threat Group-3390, ATT&CK group G0027 v1.0", "name": "Threat Group-3390 (G0027)", "domain": "mitre-enterprise", "version": "2.2", "techniques": [{"score": 1, "techniqueID": "T1203", "techniqueName": "Exploitation for Client Execution", "comment": "[Threat Group-3390](https://attack.mitre.org/groups/G0027) has exploited the Microsoft SharePoint vulnerability CVE-2019-0604."}, {"score": 1, "techniqueID": "T1210", "techniqueName": "Exploitation of Remote Services", "comment": "[Threat Group-3390](https://attack.mitre.org/groups/G0027) has exploited MS17-101 to move laterally to other systems on the network.\t"}, {"score": 1, "techniqueID": "T1100", "techniqueName": "Web Shell", "comment": "[Threat Group-3390](https://attack.mitre.org/groups/G0027) has used a variety of Web shells."}, {"score": 1, "techniqueID": "T1018", "techniqueName": "Remote System Discovery", "comment": "[Threat Group-3390](https://attack.mitre.org/groups/G0027) has used the <code>net view</code> command."}, {"score": 1, "techniqueID": "T1038", "techniqueName": "DLL Search Order Hijacking", "comment": "[Threat Group-3390](https://attack.mitre.org/groups/G0027) has performed DLL search order hijacking to execute their payload."}, {"score": 1, "techniqueID": "T1071", "techniqueName": "Standard Application Layer Protocol", "comment": "[Threat Group-3390](https://attack.mitre.org/groups/G0027) malware has used HTTP for C2."}, {"score": 1, "techniqueID": "T1050", "techniqueName": "New Service", "comment": "A [Threat Group-3390](https://attack.mitre.org/groups/G0027) tool can create a new service, naming it after the config information, to gain persistence."}, {"score": 1, "techniqueID": "T1088", "techniqueName": "Bypass User Account Control", "comment": "A [Threat Group-3390](https://attack.mitre.org/groups/G0027) tool can use a public UAC bypass method to elevate privileges."}, {"score": 1, "techniqueID": "T1055", "techniqueName": "Process Injection", "comment": "A [Threat Group-3390](https://attack.mitre.org/groups/G0027) tool can spawn svchost.exe and inject the payload into that process."}, {"score": 1, "techniqueID": "T1047", "techniqueName": "Windows Management Instrumentation", "comment": "A [Threat Group-3390](https://attack.mitre.org/groups/G0027) tool can use WMI to execute a binary."}, {"score": 1, "techniqueID": "T1060", "techniqueName": "Registry Run Keys / Startup Folder", "comment": "A [Threat Group-3390](https://attack.mitre.org/groups/G0027) tool can add the binary\u2019s path to the Registry key <code>Software\\Microsoft\\Windows\\CurrentVersion\\Run</code> to add persistence."}, {"score": 1, "techniqueID": "T1112", "techniqueName": "Modify Registry", "comment": "A [Threat Group-3390](https://attack.mitre.org/groups/G0027) tool can create a new Registry key under <code>HKEY_CURRENT_USER\\Software\\Classes\\</code>."}, {"score": 1, "techniqueID": "T1012", "techniqueName": "Query Registry", "comment": "A [Threat Group-3390](https://attack.mitre.org/groups/G0027) tool can read and decrypt stored Registry values."}, {"score": 1, "techniqueID": "T1189", "techniqueName": "Drive-by Compromise", "comment": "[Threat Group-3390](https://attack.mitre.org/groups/G0027) has extensively used strategic web compromises to target victims."}, {"score": 1, "techniqueID": "T1027", "techniqueName": "Obfuscated Files or Information", "comment": "A [Threat Group-3390](https://attack.mitre.org/groups/G0027) tool can encrypt payloads using XOR. [Threat Group-3390](https://attack.mitre.org/groups/G0027) malware is also obfuscated using Metasploit\u2019s shikata_ga_nai encoder as well as compressed with LZNT1 compression."}, {"score": 1, "techniqueID": "T1140", "techniqueName": "Deobfuscate/Decode Files or Information", "comment": "During execution, [Threat Group-3390](https://attack.mitre.org/groups/G0027) malware deobfuscates and decompresses code that was encoded with Metasploit\u2019s shikata_ga_nai encoder as well as compressed with LZNT1 compression."}, {"score": 1, "techniqueID": "T1049", "techniqueName": "System Network Connections Discovery", "comment": "[Threat Group-3390](https://attack.mitre.org/groups/G0027) has used <code>net use</code> to conduct internal discovery of systems. The group has also used quser.exe to identify existing RDP sessions on a victim."}, {"score": 1, "techniqueID": "T1016", "techniqueName": "System Network Configuration Discovery", "comment": "[Threat Group-3390](https://attack.mitre.org/groups/G0027) actors use nbtscan to discover vulnerable systems."}, {"score": 1, "techniqueID": "T1108", "techniqueName": "Redundant Access", "comment": "[Threat Group-3390](https://attack.mitre.org/groups/G0027) has deployed backup web shells and obtained OWA account credentials during intrusions that it subsequently used to attempt to regain access when evicted from a victim network."}, {"score": 1, "techniqueID": "T1028", "techniqueName": "Windows Remote Management", "comment": "[Threat Group-3390](https://attack.mitre.org/groups/G0027) has used WinRM to enable remote execution."}, {"score": 1, "techniqueID": "T1087", "techniqueName": "Account Discovery", "comment": "[Threat Group-3390](https://attack.mitre.org/groups/G0027) has used <code>net user</code> to conduct internal discovery of systems."}, {"score": 1, "techniqueID": "T1068", "techniqueName": "Exploitation for Privilege Escalation", "comment": "[Threat Group-3390](https://attack.mitre.org/groups/G0027) has used CVE-2014-6324 to escalate privileges."}, {"score": 1, "techniqueID": "T1086", "techniqueName": "PowerShell", "comment": "[Threat Group-3390](https://attack.mitre.org/groups/G0027) has used PowerShell for execution."}, {"score": 1, "techniqueID": "T1107", "techniqueName": "File Deletion", "comment": "[Threat Group-3390](https://attack.mitre.org/groups/G0027) has deleted existing logs and exfiltrated file archives from a victim."}, {"score": 1, "techniqueID": "T1126", "techniqueName": "Network Share Connection Removal", "comment": "[Threat Group-3390](https://attack.mitre.org/groups/G0027) has detached network shares after exfiltrating files, likely to evade detection."}, {"score": 1, "techniqueID": "T1059", "techniqueName": "Command-Line Interface", "comment": "[Threat Group-3390](https://attack.mitre.org/groups/G0027) has used command-line interfaces for execution."}, {"score": 1, "techniqueID": "T1089", "techniqueName": "Disabling Security Tools", "comment": "[Threat Group-3390](https://attack.mitre.org/groups/G0027) has used appcmd.exe to disable logging on a victim server."}, {"score": 1, "techniqueID": "T1119", "techniqueName": "Automated Collection", "comment": "[Threat Group-3390](https://attack.mitre.org/groups/G0027) ran a command to compile an archive of file types of interest from the victim user's directories."}, {"score": 1, "techniqueID": "T1030", "techniqueName": "Data Transfer Size Limits", "comment": "[Threat Group-3390](https://attack.mitre.org/groups/G0027) actors have split RAR files for exfiltration into parts."}, {"score": 1, "techniqueID": "T1005", "techniqueName": "Data from Local System", "comment": "[Threat Group-3390](https://attack.mitre.org/groups/G0027) ran a command to compile an archive of file types of interest from the victim user's directories."}, {"score": 1, "techniqueID": "T1133", "techniqueName": "External Remote Services", "comment": "[Threat Group-3390](https://attack.mitre.org/groups/G0027) actors look for and use VPN profiles during an operation to access the network using external VPN services."}, {"score": 1, "techniqueID": "T1105", "techniqueName": "Remote File Copy", "comment": "After re-establishing access to a victim network, [Threat Group-3390](https://attack.mitre.org/groups/G0027) actors download tools including [gsecdump](https://attack.mitre.org/software/S0008) and WCE that are staged temporarily on websites that were previously compromised but never used."}, {"score": 1, "techniqueID": "T1022", "techniqueName": "Data Encrypted", "comment": "[Threat Group-3390](https://attack.mitre.org/groups/G0027) has used RAR to compress, encrypt, and password-protect files prior to exfiltration."}, {"score": 1, "techniqueID": "T1078", "techniqueName": "Valid Accounts", "comment": "[Threat Group-3390](https://attack.mitre.org/groups/G0027) actors obtain legitimate credentials using a variety of methods and use them to further lateral movement on victim networks."}, {"score": 1, "techniqueID": "T1074", "techniqueName": "Data Staged", "comment": "[Threat Group-3390](https://attack.mitre.org/groups/G0027) has staged encrypted archives for exfiltration on Internet-facing servers that had previously been compromised with [China Chopper](https://attack.mitre.org/software/S0020)."}, {"score": 1, "techniqueID": "T1002", "techniqueName": "Data Compressed", "comment": "[Threat Group-3390](https://attack.mitre.org/groups/G0027) has used RAR to compress, encrypt, and password-protect files prior to exfiltration."}, {"score": 1, "techniqueID": "T1056", "techniqueName": "Input Capture", "comment": "[Threat Group-3390](https://attack.mitre.org/groups/G0027) actors installed a credential logger on Microsoft Exchange servers. [Threat Group-3390](https://attack.mitre.org/groups/G0027) also leveraged the reconnaissance framework, ScanBox, to capture keystrokes."}, {"score": 1, "techniqueID": "T1043", "techniqueName": "Commonly Used Port", "comment": "C2 traffic for most [Threat Group-3390](https://attack.mitre.org/groups/G0027) tools occurs over Port Numbers 53, 80, and 443."}, {"score": 1, "techniqueID": "T1073", "techniqueName": "DLL Side-Loading", "comment": "[Threat Group-3390](https://attack.mitre.org/groups/G0027) has used DLL side-loading, including by using legitimate Kaspersky antivirus variants in which the DLL acts as a stub loader that loads and executes the shell code."}, {"score": 1, "techniqueID": "T1003", "techniqueName": "Credential Dumping", "comment": "[Threat Group-3390](https://attack.mitre.org/groups/G0027) actors have used [gsecdump](https://attack.mitre.org/software/S0008) and a modified version of [Mimikatz](https://attack.mitre.org/software/S0002) called Wrapikatz to dump credentials. They have also dumped credentials from domain controllers."}, {"score": 1, "techniqueID": "T1053", "techniqueName": "Scheduled Task", "comment": "[Threat Group-3390](https://attack.mitre.org/groups/G0027) actors use [at](https://attack.mitre.org/software/S0110) to schedule tasks to run self-extracting RAR archives, which install [HTTPBrowser](https://attack.mitre.org/software/S0070) or [PlugX](https://attack.mitre.org/software/S0013) on other victims on a network."}, {"score": 1, "techniqueID": "T1046", "techniqueName": "Network Service Scanning", "comment": "[Threat Group-3390](https://attack.mitre.org/groups/G0027) actors use the Hunter tool to conduct network service discovery for vulnerable systems."}], "gradient": {"colors": ["#ffffff", "#ff6666"], "minValue": 0, "maxValue": 1}, "legendItems": [{"label": "used by Threat Group-3390", "color": "#ff6666"}]}